# deploy/av-scanner.service
[Unit]
Description=AV Scanner Service
Documentation=https://github.com/your-org/av-scanner
After=network.target clamav-daemon.service
Wants=clamav-daemon.service

[Service]
Type=simple
User=avscanner
Group=avscanner
WorkingDirectory=/opt/av-scanner

# Environment file for configuration
EnvironmentFile=-/etc/av-scanner/av-scanner.env

# Node.js application
ExecStart=/usr/bin/node /opt/av-scanner/dist/index.js

# Restart policy
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=false
ReadWritePaths=/tmp/av-scanner
ReadOnlyPaths=/var/run/clamav /var/log/clamav /var/log/ds_agent /opt/ds_agent

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=av-scanner

[Install]
WantedBy=multi-user.target
